# api/crontab
# Event Horizon periodic schedule — run by supercronic inside the `api` container.
#
# Format:  m h dom mon dow   command
# Every line dispatches a background task via `manage.py run_task <name> [key=value ...]`.
# The task's queue is auto-selected (see run_task.py); feature-gated tasks carry
# --require-flag so they stay listed here unconditionally and no-op when disabled.
#
# To change a cadence, edit this file and restart the `api` service.

# ── The pipeline ─────────────────────────────────────────────────────────
# One tick dispatches every pipeline stage that is due and has pending work
# (fetch → {analyze | annotate → refine} → aggregate → tag → route). Per-stage
# cadence, chunking, and queue live in services/stages.py — NOT here.
*/10 * * * *   python /app/manage.py run_task pipeline_tick_task

# Daily full-window aggregate sweep — the live 'aggregate' stage only clusters the
# trailing AGGREGATE_LIVE_WINDOW_HOURS (72h) each tick; this re-aggregates the full
# EVENT_STAGE_WINDOW_HOURS (168h) once daily for multi-day events past that window.
0 1 * * *      python /app/manage.py run_task aggregate_full_task

# ── Streams (light queue) ────────────────────────────────────────────────
*/5  * * * *   python /app/manage.py run_task run_stream_task name=prices      --require-flag STREAM_PRICES_ENABLED
*/15 * * * *   python /app/manage.py run_task run_stream_task name=notam      --require-flag STREAM_NOTAM_ENABLED
*/5  * * * *   python /app/manage.py run_task run_stream_task name=earthquakes --require-flag STREAM_EARTHQUAKE_ENABLED
*/15 * * * *   python /app/manage.py run_task run_stream_task name=forex       --require-flag STREAM_FOREX_ENABLED

# ── Health ───────────────────────────────────────────────────────────────
*/30 * * * *   python /app/manage.py run_task pipeline_health_task

# ── Daily LLM cron jobs ──────────────────────────────────────────────────
0 0 * * *      python /app/manage.py run_task refresh_openrouter_models_task --require-flag OPENROUTER_DYNAMIC_MODELS
0 4 * * *      python /app/manage.py run_task refresh_topics_task
0 5 * * *      python /app/manage.py run_task discover_topics_task
0 6 * * *      python /app/manage.py run_task generate_newsletter_task --require-flag NEWSLETTER_ENABLED

# ── Maintenance (light) ──────────────────────────────────────────────────
0  2 * * 0     python /app/manage.py run_task adjust_source_weights_task           --require-flag ARTICLE_IMPORTANCE_SCORING_ENABLED

# ── Forecasting ──────────────────────────────────────────────────────────
0  5  * * *    python /app/manage.py run_task train_forecast_model_task  --require-flag FORECAST_ENABLED
30 5  * * *    python /app/manage.py run_task run_forecast_task           --require-flag FORECAST_ENABLED
0 7   * * *    python /app/manage.py run_task score_forecasts_task            --require-flag FORECAST_ENABLED
0 0   * * 0    python /app/manage.py run_task backfill_prices_task            --require-flag FORECAST_ENABLED
